# base/supervisor/supervisor-x11vnc.ini
[program:x11vnc]
command=/usr/bin/x11vnc -display "%(ENV_DISPLAY)s" -xkb -noxrecord -noxfixes -noxdamage -forever -nopw
redirect_stderr=true
stdout_logfile=/var/log/x11vnc/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/x11vnc/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
