# supervisor/supervisor-xvfb.ini
[program:xvfb]
command=/usr/bin/Xvfb "%(ENV_DISPLAY)s" -screen "%(ENV_RESOLUTION)s" -ac
redirect_stderr=true
stdout_logfile=/var/log/xvfb/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/xvfb/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
