# supervisor/supervisor-novnc.ini
[program:novnc]
command=/usr/bin/novnc_server --web /usr/share/novnc/
redirect_stderr=true
stdout_logfile=/var/log/novnc/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/novnc/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
