# supervisor/supervisor-xfce4.ini
[program:xfce4]
command=/usr/bin/startxfce4
redirect_stderr=true
stdout_logfile=/var/log/xfce4/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/xfce4/stderr
stderr_logfile_maxbytes=5MB
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
